# systemd-unit-files/disk-keepalive.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /disk-keepalive.tmp
